# custom/includes.chroot/etc/update-motd.d/10-uname
#!/bin/sh

printf "\n"
printf "Welcome to %s (%s %s %s)\n" "$(uname -n)" "$(uname -o)" "$(uname -r)" "$(uname -m)"
